# MongooseIm_V1.5/MongooseIM/deb-scripts/postinst
#!/bin/sh
# postinst script for mongooseim
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

NAME="mongooseim"
PKG_ROOT="/usr/lib/${NAME}/"

. /usr/share/debconf/confmodule

setup()
{
	if ! dpkg-statoverride --list $PKG_ROOT >/dev/null; then
		chown -R $NAME:$NAME $PKG_ROOT
	fi

        cp $PKG_ROOT/bin/mongooseim /etc/init.d/mongooseim
        cp $PKG_ROOT/bin/mongooseimctl /usr/sbin/mongooseimctl
	sed -i -e "s*\(RUNNER_SCRIPT_DIR=\).\**\1${PKG_ROOT}bin*" /usr/sbin/mongooseimctl
        sed -i -e "s*\(RUNNER_SCRIPT_DIR=\).\**\1${PKG_ROOT}bin*" /etc/init.d/mongooseim
	sed -i -e "s*\(export EJABBERD_CONFIG_PATH\).\**\1\nexport HOME=/usr/lib/mongooseim*" /etc/init.d/mongooseim
}

case "$1" in
    configure)
	adduser --quiet --system --shell /bin/sh --group --home $PKG_ROOT $NAME 2>/dev/null
	setup "$@"
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit
echo "installing service script..."
if [ -x "/etc/init.d/mongooseim" ]; then
	if [ ! -e "/etc/init/mongooseim.conf" ]; then
		update-rc.d mongooseim defaults >/dev/null
	fi
        PID=`ps aux | grep progname| grep mongooseim | grep -v grep | awk '{print $2}'`
	if [ "$PID" = "" ]; then
           echo "Start"
	   invoke-rc.d mongooseim start || exit $?
        else
	   echo "Restart"
           invoke-rc.d mongooseim stop || exit $?
           invoke-rc.d mongooseim start || exit $?
        fi
fi
# End automatically added section


exit 0
